# cluster-sync/global/etc/profile.d/userdata.sh
# Application binaries on the shared filesystem
export PATH="/userdata/apps/bin:$PATH"
